# private_dot_config/systemd/user/sing-box.service
[Unit]
Description=sing-box proxy service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/sing-box run -c %h/.config/sing-box/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
