# private_dot_config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper from Unsplash

[Service]
Type=oneshot
ExecStart=%h/.config/sway/scripts/wallpaper.sh
